# crontab
0 0 * * * /path/to/your/python/interpereter /path/to/your/python/script
